# cleanup/trigger-cleanup.service
[Unit]
Description=Trigger automated cleanup task
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker container prune --force
TimeoutStartSec=5min

[Install]
WantedBy=multi-user.target
